# deploy/systemd/seraya-website.service
[Unit]
Description=PT Sejahtera Raya Grup - Corporate Website
After=network.target

[Service]
Type=simple
User=fajar
Group=fajar
WorkingDirectory=/home/fajar/seraya/SerayaGroup
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=3000
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
